# deployment/apache-answers-bot.service
[Unit]
Description=Apache Answers - Microsoft Teams Integration
Documentation=https://github.com/your-repo
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=apache-answers-bot
Group=apache-answers-bot
WorkingDirectory={{APP_DIR}}
Environment=NODE_ENV=production
EnvironmentFile={{APP_DIR}}/.env
ExecStart=/usr/bin/node {{APP_DIR}}/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=apache-answers-bot

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/apache-answers-bot

# Resource limits (optional - adjust as needed)
LimitNOFILE=65536
# MemoryMax=512M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
